# etc-trash/home/config/systemd/user/suspend.target
[Unit]
Description=User level suspend target
StopWhenUnneeded=yes
Wants=user-hibernate.service
