# etc-trash/home/config/systemd/user/user-hibernate.service
[Unit]
Description=Disable lid switch on hibernate start

[Service]
ExecStart=

[Install]
Wants=suspend.target
